# .config/autostart/lxpolkit.desktop
[Desktop Entry]
Name=Lxpolkit
Exec=/usr/bin/lxpolkit
Terminal=false
Type=Application
NoDisplay=true
